# modules/lan/files/etc/uci-defaults/89-lan

uci batch << EOF

set network.lan.gateway='$LAN_GATEWAY'
set network.lan.ipaddr='$LAN_IP'

commit
EOF
