# modules/add-all-device-to-lan/files/etc/uci-defaults/99-add-all-device-to-lan
wan_device=$(uci get network.wan.device)
all_devices=$(ip -o link show | awk -F': ' '{print $2}' | grep "^eth")

# 删除 wan_device
all_devices=$(echo "$all_devices" | tr ' ' '\n' | grep -v "^$wan_device$" | xargs)

# 清空旧端口
uci -q delete network.lan.ports

# 重新添加剩余设备
for device in $all_devices; do
    uci add_list network.lan.ports=$device
done

uci commit network
/etc/init.d/network restart
